# signals_box_ctl/service_ctl.service
# https://blog.miguelgrinberg.com/post/running-a-flask-application-as-a-service-with-systemd
[Unit]
Description=RF Service Manager
After=network.target

[Service]
WorkingDirectory=/opt/signals_box_ctl
ExecStart=/opt/signals_box_ctl/signals/bin/gunicorn -b 0.0.0.0:8080 -w 4 wsgi:app
Environment="PATH=/opt/signals_box_ctl/signals/bin/"
Restart=always

[Install]
WantedBy=multi-user.target
